# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:NetThermPrintServer]
command=/usr/bin/python3 /server.py
priority=30

